# redissentinel/redisproxy/supervisord.conf
[unix_http_server]
file=/data/www/redisproxy/supervisor_redisproxy_PROXY_PORT.sock

[supervisord]
logfile=/data/logs/redisproxy/supervisor_redisproxy_PROXY_PORT.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=100MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/data/logs/redisproxy/supervisor_redisproxy_PROXY_PORT.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///data/www/redisproxy/supervisor_redisproxy.sock 

[program:redisproxy]
command=/data/www/redisproxy/bin/redisproxy -http_port=PROXY_HTTPPORT -port=PROXY_PORT -idc=IDC -namespace=NAMESPACE -zk_addrs=ZK_ADDRS
autostart=true
autorestart=true
numprocs=1                    
directory=/data/www/redisproxy/             
stdout_logfile=/data/logs/redisproxy/redisproxy_out_PROXY_PORT.log     
stderr_logfile=/data/logs/redisproxy/redisproxy_err_PROXY_PORT.log  

[eventlistener:eventexec]
command=eventexec -e "ps -ef|grep supervisord|grep PROXY_PORT|grep -v grep|awk '{print $2}'|xargs kill"
events=PROCESS_STATE_FATAL
